# pi-server.service
[Unit]
Description=Pi Server (Express + Socket.io)
After=network.target

[Service]
Type=simple
User=zerocool
WorkingDirectory=/home/zerocool/pi-server
ExecStart=/usr/bin/node /home/zerocool/pi-server/server.js
Restart=on-failure
RestartSec=5
Environment=PORT=3000
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
